# Water_Atm_Firmware_Uploader_and_Updater/python_init
#!/bin/sh /etc/rc.common

START=99

start(){
	python /root/main.py & python /root/Firmware_Updater.py
}

# Rango is back again

